# odoo-scripts/crontab
# m h dom mon dow user  command
0 1 * * * postgres /var/scripts/dump_db.sh
0 5 * * * postgres /var/scripts/housekeeping.sh
